# usb_udev_check_states/test_USB_state.rules
RUN+="/home/icewindful/udev_event_script"

